# cron/crontab
*  *  *  *  * cron date >> /tmp/cron.log
